# systemd/bc250-health.timer
[Unit]
Description=BC250 Service Health Check Timer (every 15 min)

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
AccuracySec=30s

[Install]
WantedBy=timers.target
