# systemd/ollama.service.d/override.conf
[Service]
Environment="OLLAMA_KEEP_ALIVE=30m"
Environment="OLLAMA_MAX_LOADED_MODELS=1"
Environment="OLLAMA_VULKAN=1"
Environment="OLLAMA_FLASH_ATTENTION=1"
Environment="OLLAMA_GPU_OVERHEAD=0"
Environment="OLLAMA_CONTEXT_LENGTH=24576"
Environment="OLLAMA_MAX_QUEUE=4"
OOMScoreAdjust=-1000
